# puppetmaster5/supervisor/autocommit.ini
[program:autocommit]

command=/opt/utils/autocommit/autocommit.sh -r /etc/puppetlabs/puppet/.repo/ssl-repo -p -f /etc/puppetlabs/puppet/.repo/ssl-repo/.autocommit.enabled
autostart=true
autorestart=true
