# puppetmaster5/supervisor/crond.ini
[program:crond]

command=/usr/sbin/crond -n -m off
autostart=true
autorestart=true
