# puppetmaster5/supervisor/puppetmaster.ini
[program:puppetmaster]

command=/usr/local/bin/start_puppetmaster.sh
autostart=true
autorestart=true
